# extras/systemd/isopod-connected.service
# While isopod.service is resilient to the loss of its rsync target, it would
# be a bad experience for an isopod terminal user to rip a bunch of discs
# before realizing they're not going anywhere. This check is designed to run
# once per system boot, so that isopod.service can activate its OnFailure=
# handling for unrecoverable mistakes in the terminal configuration (e.g. you
# typo'd the WiFi password before shipping it to the end user).

[Unit]
Description=Check connectivity to the isopod target
After=network-online.target
StartLimitIntervalSec=6m
StartLimitBurst=60

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/default/isopod
ExecStart=rsync --list-only ${ISOPOD_TARGET}
StandardOutput=null
StandardError=journal
TimeoutStartSec=60s
Restart=on-failure
RestartSec=3s
RestartPreventExitStatus=20 SIGTERM SIGKILL

[Install]
WantedBy=network-online.target
