# extras/systemd/isopod-systemfail.timer
[Unit]
Description=Report failure on the E-Ink display after a delay

[Timer]
OnActiveSec=30s
